# mission2/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/etc/init.d/postfix start

[program:dovecot]
command=/etc/init.d/dovecot start